# shadowsocks.conf
[program:shadowsocks]
command=ssserver -c /opt/shadowsocks/config.json
autostart=true
autorestart=true 
user=nobody